# user/system_service_example/etc/systemd/system/networking-check.service
[Unit]
After=snmpd.service

[Service]
ExecStart=/usr/local/bin/networking-check.sh

[Install]
WantedBy=default.target
