# EXPC-1519/romfs_debian9/etc/udev/rules.d/01-rename_uart_naming.rules
# The UART port order is reserve, we use this rule to swap /dev/ttyS0 and /dev/ttyS1
KERNEL=="ttyS0", RUN+="/bin/bash -c '/bin/mv /dev/ttyS0 /dev/ttyStmp0 ; /bin/mv /dev/ttyS1 /dev/ttyS0; /bin/mv /dev/ttyStmp0 /dev/ttyS1'"
